# deploy/systemd/yclients-api@.service
[Unit]
Description=YClients BI FastAPI service
After=network.target docker.service
Wants=network.target

[Service]
Type=simple
User=%i
WorkingDirectory=/opt/yclients_bi_system
Environment=PYTHONUNBUFFERED=1
ExecStart=/opt/yclients_bi_system/.venv/bin/python api.py
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
