# deploy/vm/yclients-auto-deploy.service
[Unit]
Description=Deploy YClients BI from GitHub
After=network-online.target docker.service
Wants=network-online.target docker.service

[Service]
Type=oneshot
Environment=APP_DIR=/opt/yclients_bi_system
Environment=BRANCH=main
ExecStart=/opt/yclients_bi_system/deploy/vm/yclients-auto-deploy.sh
