# deploy/vm/yclients-auto-deploy.timer
[Unit]
Description=Check GitHub for YClients BI deploys every 5 minutes

[Timer]
OnBootSec=2min
OnUnitActiveSec=5min
Persistent=true
RandomizedDelaySec=30s
Unit=yclients-auto-deploy.service

[Install]
WantedBy=timers.target
